# demo/files/appwatcher.path
[Path]
Unit=appwatcher.service
PathModified=/opt/gateway

[Install]
WantedBy=multi-user.target
